# php-fpm-nginx/supervisord.conf
[supervisord]
nodaemon=true

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock
[program:php-fpm]
command=php-fpm
[program:nginx]
command=nginx -g "daemon off;"

